# dist/supervisor/service-websocket.conf
[program:service-websocket]
command=nodejs /usr/share/yoctu/service-websocket/app.js
directory=/usr/share/yoctu/service-websocket
autostart=true
autorestart=true
stopwaitsecs=5
startretries=10
exitcodes=0,1,2,10,11
stderr_logfile=/var/log/supervisor/service-websocket-stderr.log
stdout_logfile=/var/log/supervisor/service-websocket-stdout.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
user=root
